# .config/autostart/Flameshot.desktop
[Desktop Entry]
Exec=/bin/bash -c "sleep 10; flameshot &"